# cron/analyze.cron
0 */23 * * * /www/sites/ala-poker.com/inc/analyze.php >> analyze.run
